# navra-server/systemd/navra.socket
[Unit]
Description=navra — MCP Server Socket

[Socket]
ListenStream=%t/navra/navra.sock
SocketMode=0600
DirectoryMode=0700

[Install]
WantedBy=sockets.target
