# mcpd-docs/systemd/mcpd-docs.service
[Unit]
Description=mcpd-docs — MCP Document Server
Documentation=https://github.com/user/mcpd
After=default.target

[Service]
Type=notify
ExecStart=%h/.cargo/bin/mcpd-docs serve
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/mcpd-docs
RuntimeDirectory=mcpd
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

[Install]
WantedBy=default.target
